# supervisord.conf
[program:cpu]
command=python /tmp/cpuload/cpuloader.py
